# social-platform/deploy/watchparty.service.example
[Unit]
Description=Watchparty social platform
After=network.target

[Service]
Type=simple
User=watchparty
Group=watchparty
WorkingDirectory=/opt/watchparty/social-platform
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3000
Environment=TRUST_PROXY=1
Environment=DATA_DIR=/var/lib/watchparty
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
